# launchagents/com.openclaw.graphiti-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.graphiti-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- UPDATE THIS PATH to match your setup -->
        <string>/Users/YOUR_USERNAME/clawd/scripts/graphiti-sync-sessions.py</string>
    </array>
    <key>StartInterval</key>
    <!-- Runs every 30 minutes (1800 seconds) -->
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <!-- UPDATE THIS PATH -->
    <string>/Users/YOUR_USERNAME/.clawdbot/logs/graphiti-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.clawdbot/logs/graphiti-sync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <!-- UPDATE THIS -->
        <string>/Users/YOUR_USERNAME</string>
    </dict>
</dict>
</plist>
